# crontab
0 9 * * 1-5 python /main.py > /dev/stdout
